# sec_audit.crontab
PATH=/usr/local/bin:/bin:/usr/bin:/usr/sbin
5 */1 * * * /root/bin/compare-sg.sh
